# initrd/source/etc/udev/rules.d/99-tmn.rules
ACTION=="add", KERNEL=="ttyUSB2", RUN+="/bin/sh -c 'killall -INT pppd; sleep 2; pppd call provider'" 
